# deploy/aio/s6-rc.d/migrate/up
/command/with-contenv /usr/local/bin/tessio-migrate
