# deploy/aio/s6-rc.d/worker/run
#!/command/with-contenv sh
cd /app/worker
exec s6-setuidgid node node dist/index.js
