# deploy/aio/s6-rc.d/runner/run
#!/command/with-contenv sh
cd /app/runner
exec s6-setuidgid node node dist/index.js
